# systemd/agent-oci-to-umbrella.service
[Unit]
Description=OCI to Umbrella BYOD Transfer Agent
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=oci-umbrella-agent
Group=oci-umbrella-agent
WorkingDirectory=/opt/agent-oci-to-umbrella
ExecStart=/usr/local/bin/agent-oci-to-umbrella run --config /etc/agent-oci-to-umbrella/config.yaml
Restart=on-failure
RestartSec=30
StandardOutput=journal
StandardError=journal
SyslogIdentifier=agent-oci-to-umbrella

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/agent-oci-to-umbrella/logs /opt/agent-oci-to-umbrella/state

# Environment
Environment="PATH=/usr/local/bin:/usr/bin:/bin"

[Install]
WantedBy=multi-user.target
