# systemd/oracle-focus-agent.service
[Unit]
Description=Oracle FOCUS to Umbrella BYOD Transfer Agent
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=oracle-agent
Group=oracle-agent
WorkingDirectory=/opt/oracle-focus-agent
ExecStart=/usr/local/bin/oracle-focus-agent run --config /etc/oracle-focus-agent/config.yaml
Restart=on-failure
RestartSec=30
StandardOutput=journal
StandardError=journal
SyslogIdentifier=oracle-focus-agent

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/oracle-focus-agent/logs /opt/oracle-focus-agent/state

# Environment
Environment="PATH=/usr/local/bin:/usr/bin:/bin"

[Install]
WantedBy=multi-user.target
